# aports/device/testing/device-oneplus-hotdog/device-oneplus-hotdog.post-install
#!/bin/sh

# postmarketos-mkinitfs edge expects this path while Alpine losetup is
# still installed in /sbin in this rootfs.
mkdir -p /usr/sbin
ln -sf /sbin/losetup /usr/sbin/losetup

# mdev does not always create the Android partition nodes/symlinks early enough
# for SSH-driven flash tests. Derive them from sysfs on boot.
mkdir -p /etc/local.d
cat > /etc/local.d/hotdog-devnodes.start <<'EOF'
#!/bin/sh
set -eu

mkdir -p /dev/disk/by-partlabel /dev/disk/by-name /dev/block/by-name /dev/block

for uevent in /sys/class/block/*/uevent; do
	[ -r "$uevent" ] || continue
	devname=""
	partname=""
	major=""
	minor=""
	while IFS='=' read -r key value; do
		case "$key" in
			DEVNAME) devname="$value" ;;
			PARTNAME) partname="$value" ;;
			MAJOR) major="$value" ;;
			MINOR) minor="$value" ;;
		esac
	done < "$uevent"

	[ -n "$devname" ] || continue
	if [ ! -b "/dev/$devname" ] && [ -n "$major" ] && [ -n "$minor" ]; then
		mknod "/dev/$devname" b "$major" "$minor" 2>/dev/null || true
	fi
	[ -b "/dev/$devname" ] || continue

	ln -sf "../$devname" "/dev/block/$devname" 2>/dev/null || true
	[ -n "$partname" ] || continue
	ln -sf "../../$devname" "/dev/disk/by-partlabel/$partname" 2>/dev/null || true
	ln -sf "../../$devname" "/dev/disk/by-name/$partname" 2>/dev/null || true
	ln -sf "../../$devname" "/dev/block/by-name/$partname" 2>/dev/null || true
done
EOF
chmod 0755 /etc/local.d/hotdog-devnodes.start
rc-update add local default 2>/dev/null || true
